# deploy/swat-ws.service
[Unit]
Description=Swat WebSocket relay
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/swat/server
Environment=NODE_ENV=production
Environment=PORT=8080
ExecStart=/usr/bin/node index.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
